# usr/etc/skel/.config/autostart/rslsync.desktop
[Desktop Entry]
Exec=~/.nix-profile/bin/rslsync --config ~/.config/rslsync/rslsync.conf
Icon=
Name=rslsync
Path=
Terminal=False
Type=Application
